# L3akCTF2025/magical-oracle/build/supervisord.conf
; supervisord.conf
[supervisord]
nodaemon=true

[program:socat]
command=/opt/app/run.sh
user=root
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
