# HackDearbornCTF-RewindReality/web/retro-bot/build/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python3 /app/app.py
autostart=true
autorestart=true