# HackDearbornCTF-RewindReality/web/time-will-tell/dist/supervisord.conf
[supervisord]
nodaemon=true

[program:timing-ctf]
command=python3 app.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/timing-ctf.err.log
stdout_logfile=/var/log/timing-ctf.out.log