# HackDearbornCTF-RewindReality/web/retro-game-search/build/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python3 app.py
user=root
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0