# HackDearbornCTF-RewindReality/web/time-will-tell/build/supervisord.conf
[supervisord]
nodaemon=true

[program:timing-ctf]
command=python3 app.py
directory=/app
autostart=true
autorestart=true